# api-list/scheduler/airilab-worker.service
[Unit]
Description=AiriLab API Worker Daemon
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/.openclaw/skills/api-list/scheduler
ExecStart=/usr/bin/python3 /home/ec2-user/.openclaw/skills/api-list/scheduler/worker.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airilab-worker

# 环境变量
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
